# toolchain/rpcd.service
[Unit]
Description=rpcd
After=ubus.service

[Service]
ExecStart=/usr/sbin/rpcd
Restart=on-failure

[Install]
WantedBy=multi-user.target
